# deploy/vpn-bot.user.service
# Пользовательский systemd (без sudo для vpn-bot после первичной настройки root).
# Каталог на сервере: ~/.config/systemd/user/vpn-bot.service
#
# Один раз под root:
#   loginctl enable-linger vpn-bot
#   install -d -o vpn-bot -g vpn-bot /home/vpn-bot/.config/systemd/user
#   cp deploy/vpn-bot.user.service /home/vpn-bot/.config/systemd/user/vpn-bot.service
#   chown vpn-bot:vpn-bot /home/vpn-bot/.config/systemd/user/vpn-bot.service
#
# Под vpn-bot:
#   systemctl --user daemon-reload
#   systemctl --user enable --now vpn-bot.service
#
# Логи: journalctl --user -u vpn-bot.service -f

[Unit]
Description=Telegram VPN Bot (AccessGranted / 3X-UI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/vpn-bot/apps/vpn-bot
ExecStart=/home/vpn-bot/apps/vpn-bot/.venv/bin/python -m src.app
Restart=on-failure
RestartSec=10
KillMode=mixed
TimeoutStopSec=30

[Install]
WantedBy=default.target
