# deploy/vpn-bot.service
# systemd unit для этого репозитория на вашем VPS.
# Размещение на сервере (как у вас уже настроено):
#   пользователь: vpn-bot
#   каталог:     /home/vpn-bot/apps/vpn-bot  (корень репозитория, есть src/)
#   venv:        /home/vpn-bot/apps/vpn-bot/.venv
#   .env:        /home/vpn-bot/apps/vpn-bot/.env
#
# Установка (от root, см. deploy/README.md):
#   cp deploy/vpn-bot.service /etc/systemd/system/vpn-bot.service
#   systemctl daemon-reload
#   systemctl enable --now vpn-bot.service

[Unit]
Description=Telegram VPN Bot (AccessGranted / 3X-UI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vpn-bot
Group=vpn-bot
WorkingDirectory=/home/vpn-bot/apps/vpn-bot
ExecStart=/home/vpn-bot/apps/vpn-bot/.venv/bin/python -m src.app
Restart=on-failure
RestartSec=10
KillMode=mixed
TimeoutStopSec=30

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
